# etc/cron.d/splat-plot
# Every 5 minutes, process pending splat plots...
# Our script is rate limited (see private/splat.ini)
*/5 * * * *	www-remote /svc/www/remote/bin/splat-cron.php
